# freebsd-ports-diablo/graphics/mesa-dev/pkg-plist
bin/mesa-overlay-control.py
etc/OpenCL/vendors/mesa.icd
include/EGL/egl.h
include/EGL/eglext.h
include/EGL/eglextchromium.h
include/EGL/eglmesaext.h
include/EGL/eglplatform.h
include/GL/gl.h
include/GL/glcorearb.h
include/GL/glext.h
include/GL/glx.h
include/GL/glxext.h
include/GL/internal/dri_interface.h
include/GL/osmesa.h
include/GLES/egl.h
include/GLES/gl.h
include/GLES/glext.h
include/GLES/glplatform.h
include/GLES2/gl2.h
include/GLES2/gl2ext.h
include/GLES2/gl2platform.h
include/GLES3/gl3.h
include/GLES3/gl31.h
include/GLES3/gl32.h
include/GLES3/gl3ext.h
include/GLES3/gl3platform.h
include/KHR/khrplatform.h
include/d3dadapter/d3dadapter9.h
include/d3dadapter/drm.h
include/d3dadapter/present.h
include/gbm.h
%%INTEL%%include/vulkan/vulkan_intel.h
include/xa_composite.h
include/xa_context.h
include/xa_tracker.h
lib/d3d/d3dadapter9.so
lib/d3d/d3dadapter9.so.1
lib/d3d/d3dadapter9.so.1.0.0
%%INTEL%%lib/dri/i915_dri.so
%%INTEL%%lib/dri/i965_dri.so
%%IRIS%%lib/dri/iris_dri.so
lib/dri/kms_swrast_dri.so
lib/dri/nouveau_dri.so
lib/dri/nouveau_drv_video.so
%%LEGACY%%lib/dri/nouveau_vieux_dri.so
%%LEGACY%%lib/dri/r200_dri.so
lib/dri/r300_dri.so
lib/dri/r600_dri.so
lib/dri/r600_drv_video.so
%%LEGACY%%lib/dri/radeon_dri.so
lib/dri/radeonsi_dri.so
lib/dri/radeonsi_drv_video.so
lib/dri/swrast_dri.so
lib/dri/virtio_gpu_dri.so
lib/dri/vmwgfx_dri.so
lib/gallium-pipe/pipe_nouveau.so
lib/gallium-pipe/pipe_r300.so
lib/gallium-pipe/pipe_r600.so
lib/gallium-pipe/pipe_radeonsi.so
lib/gallium-pipe/pipe_swrast.so
lib/gallium-pipe/pipe_vmwgfx.so
lib/libEGL.so
lib/libEGL.so.1
lib/libEGL.so.1.0.0
lib/libGL.so
lib/libGL.so.1
lib/libGL.so.1.2.0
lib/libGLESv1_CM.so
lib/libGLESv1_CM.so.1
lib/libGLESv1_CM.so.1.1.0
lib/libGLESv2.so
lib/libGLESv2.so.2
lib/libGLESv2.so.2.0.0
lib/libMesaOpenCL.so
lib/libMesaOpenCL.so.1
lib/libMesaOpenCL.so.1.0.0
lib/libOSMesa.so
lib/libOSMesa.so.8
lib/libOSMesa.so.8.0.0
lib/libVkLayer_MESA_overlay.so
lib/libXvMCnouveau.so
lib/libXvMCr600.so
lib/libgbm.so
lib/libgbm.so.1
lib/libgbm.so.1.0.0
lib/libglapi.so
lib/libglapi.so.0
lib/libglapi.so.0.0.0
%%SWR%%lib/libswrAVX.so
%%SWR%%lib/libswrAVX.so.0
%%SWR%%lib/libswrAVX.so.0.0.0
%%SWR%%lib/libswrAVX2.so
%%SWR%%lib/libswrAVX2.so.0
%%SWR%%lib/libswrAVX2.so.0.0.0
%%INTEL%%lib/libvulkan_intel.so
lib/libvulkan_radeon.so
lib/libxatracker.so
lib/libxatracker.so.2
lib/libxatracker.so.2.5.0
lib/vdpau/libvdpau_nouveau.so
lib/vdpau/libvdpau_nouveau.so.1
lib/vdpau/libvdpau_nouveau.so.1.0
lib/vdpau/libvdpau_nouveau.so.1.0.0
lib/vdpau/libvdpau_r300.so
lib/vdpau/libvdpau_r300.so.1
lib/vdpau/libvdpau_r300.so.1.0
lib/vdpau/libvdpau_r300.so.1.0.0
lib/vdpau/libvdpau_r600.so
lib/vdpau/libvdpau_r600.so.1
lib/vdpau/libvdpau_r600.so.1.0
lib/vdpau/libvdpau_r600.so.1.0.0
lib/vdpau/libvdpau_radeonsi.so
lib/vdpau/libvdpau_radeonsi.so.1
lib/vdpau/libvdpau_radeonsi.so.1.0
lib/vdpau/libvdpau_radeonsi.so.1.0.0
libdata/pkgconfig/d3d.pc
libdata/pkgconfig/dri.pc
libdata/pkgconfig/egl.pc
libdata/pkgconfig/gbm.pc
libdata/pkgconfig/gl.pc
libdata/pkgconfig/glesv1_cm.pc
libdata/pkgconfig/glesv2.pc
libdata/pkgconfig/osmesa.pc
libdata/pkgconfig/xatracker.pc
share/drirc.d/00-mesa-defaults.conf
share/vulkan/explicit_layer.d/VkLayer_MESA_overlay.json
%%INTEL%%share/vulkan/icd.d/intel_icd.x86_64.json
share/vulkan/icd.d/radeon_icd.%%VK_ARCH%%.json
